# board/raspberrypi-cctv/rootfs-overlay/etc/profile.d/term.sh
#!/usr/bin/bash

export TERM=vt102
[ ! -x /usr/bin/resize ] || /usr/bin/resize 1>/dev/null
